# etc/cron.d/zabbix-nfs-stats
* * * * * root /usr/bin/zabbix-nfs-stats > /dev/null
